# reel-to-recipe.service
[Unit]
Description=Reel to Recipe MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/projects/tiktok-recipe
ExecStart=%h/projects/tiktok-recipe/.venv/bin/python mcp_server.py
Restart=on-failure
RestartSec=5
Environment=WHISPER_MODEL=base
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
